# .docker/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=warn
user=root
# CRITICAL FOR HIGH CONCURRENCY:
# Default is 1024. With 700 users + DB connections, you might hit the limit.
minfds=10000

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# --- OCTANE: MAX THROUGHPUT MODE ---
[program:octane]
# WORKERS=8:
# We have ~1.8 vCPU per pod.
# Since your app is I/O bound (waiting on DB), we can run ~4 workers per vCPU.
# This raises concurrent processing from 2 requests -> 8 requests per pod.
#
# MAX-REQUESTS=500:
# Lowered slightly to recycle memory faster since we have more workers now.
command=php -d opcache.enable_cli=1 -d opcache.jit=1255 /var/www/html/artisan octane:frankenphp --host=0.0.0.0 --port=8000 --workers=5 --max-requests=500
user=appuser
autostart=true
autorestart=true
priority=1
stopsignal=SIGTERM
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- QUEUE: BACKGROUND PRIORITY ---
[program:queue-worker]
# We keep this at 1. We want the CPU focused on the USER (Octane), not background jobs.
# If queues pile up, Autoscaling will create more pods, giving you more queue workers automatically.
command=php /var/www/html/artisan queue:work --tries=3 --timeout=120 --max-jobs=500 --max-time=3600
process_name=%(program_name)s_%(process_num)02d
numprocs=1
user=appuser
autostart=true
autorestart=true
priority=2
stopsignal=SIGTERM
stopwaitsecs=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# --- SCHEDULER: MATCHMAKING & PERIODIC TASKS ---
# Runs on every pod but onOneServer() in the schedule ensures only one
# pod actually executes each scheduled job (requires CACHE_STORE=redis).
[program:scheduler]
command=php /var/www/html/artisan schedule:work
user=appuser
autostart=true
autorestart=true
priority=3
stopsignal=SIGTERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
